# docker/cron/crontab
0 3 * * * cd /app && python /app/src/pipelines/run_etl.py >> /var/log/cron.log 2>&1
